# deploy/supervisor/tbn-mlm-workers.conf
; Supervisor — TBN Living MLM API (Ubuntu 24 / VPS)
; Instalar: sudo apt install supervisor
; Copiar a: /etc/supervisor/conf.d/tbn-mlm-workers.conf
; Ajustar paths y usuario www-data

[program:tbn-mlm-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tbn-api/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600 --timeout=120 --queue=high,default,mail,binary,residual,withdrawals,low
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/tbn-api/storage/logs/worker.log
stopwaitsecs=3600

[program:tbn-mlm-scheduler]
process_name=%(program_name)s
command=php /var/www/tbn-api/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/tbn-api/storage/logs/scheduler.log
